# kiosk-script/systemd/org.gnome.Kiosk.Script.service.in
[Unit]
Description=Kiosk script
BindsTo=gnome-session.target
After=gnome-session.target

[Service]
ExecStart=@bindir@/gnome-kiosk-script
Restart=always
