# kiosk-script/systemd/org.gnome.Kiosk.Script.service.in
[Unit]
Description=Kiosk script
BindsTo=gnome-session.target
After=gnome-session.target
OnFailure=gnome-session-shutdown.target
OnSuccess=gnome-session-shutdown.target

[Service]
ExecStart=@bindir@/gnome-kiosk-script
ExitType=cgroup
